# package/thirdparty/xbmc/system.d/mali-load.service
[Unit]
Description=Load Mali/Ump modules

[Service]
Type=oneshot
ExecStart=/lib/amlinux/mali-load restart
RemainAfterExit=yes
StartLimitInterval=0

[Install]
WantedBy=xbmc.service
