# package/thirdparty/xbmc/system.d/mount-volumes.service
[Unit]
Description=Mount user volumes
DefaultDependencies=false
Before=local-fs-pre.target

[Service]
Type=oneshot
ExecStart=/lib/amlinux/mount-volumes
RemainAfterExit=yes
StartLimitInterval=0

[Install]
WantedBy=local-fs-pre.target
