# package/thirdparty/xbmc/system.d/xbmc-halt.service
[Unit]
Description=XBMC halt script
After=xbmc.service
Before=halt.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=-/bin/sh /lib/amlinux/xbmc-halt
RemainAfterExit=yes

[Install]
WantedBy=halt.target
