# package/wpa_supplicant/system.d/wpa_supplicant.service
[Unit]
Description=WPA supplicant

[Service]
Type=forking
ExecStart=/usr/sbin/wpa_supplicant -u -B -iwlan0

[Install]
WantedBy=multi-user.target
